# roles/openvpn/templates/openvpn.cron
# this file executes xs-vpn hourly to test for tunnel, and restart as required
25 *  *  *  * root /usr/bin/xs-vpn > /dev/null
