# altfiles/etc/systemd/system/moodle-xs.service.in
[Unit]
Description=Moodle Course Management OLPC edition
after=httpd.service

[Service]
type=oneshot
ExecStart=/etc/sysconfig/olpc-scripts/moodle-xs-init start

[Install]
WantedBy=multi-user.target
