# altfiles/etc/systemd/system/moodle-xs.service.in
[Unit]
Description=Moodle Course Management OLPC edition

[Service]
ExecStart=/etc/sysconfig/olpc-scripts/moodle-xs-init

[Install]
WantedBy=multi-user.target
