# roles/postgresql/templates/postgresql-xs.service
# --postgress_xs.service  # ssystemd start up file modifications
.include /lib/systemd/system/postgresql.service
[Service]
Environment=PGDATA=/library/pgsql-xs
