# plugins.d/ejabberd/files/ejabberd-xs.service.in
# --ejabberd_xs.service  
# systemd start up file modifications
.include /usr/lib/systemd/system/ejabberd.service
[Service]
EnvironmentFile=/etc/sysconfig/ejabberd-xs
