# deprecated/plugins.d/stats/files/sugar-stats-server.service.in
[Unit]
Description=Provides the sugar-stats collection server 
After=syslog.target network.target

[Service]
Type=forking
PIDFile=/library/sugar-stats/sugar-stat.pid
ExecStart=/etc/init.d/sugar-stats-server start 
ExecStop=/etc/init.d/sugar-stats-server stop

[Install]
WantedBy=multi-user.target
