# plugins.d/vnc/vnc.manifest
/usr/share/xs-config/cfg/etc/httpd/conf.d/novnc.conf
/usr/share/xs-config/cfg/etc/sysconfig/vnc/self.pem
/usr/share/xs-config/cfg/etc/sysconfig/vnc/xfce4-panel.xml
/usr/share/xs-config/cfg/etc/sysconfig/vnc/xstartup
/usr/share/xs-config/cfg/etc/systemd/system/vncserver@.service
/usr/share/xs-config/cfg/etc/systemd/system/vncserverweb.service
/usr/share/xs-config/cfg/html/top/es/vnc_auto.html
/usr/share/xs-config/cfg/html/top/vnc
/usr/share/xs-config/cfg/html/top/vnc.js
