# plugins.d/postgresql/postgresql.manifest
/etc/idmgr.conf.in
/etc/pgsql-xs/postgresql-1024m.conf
/etc/pgsql-xs/postgresql-2048m.conf
/etc/pgsql-xs/postgresql-256m.conf
/etc/pgsql-xs/postgresql-512m.conf
/library/pgsql-xs/
/etc/sysconfig/pgsql/pgsql-xs
/usr/share/xs-config/cfg/etc/systemd/system/postgresql-xs.service
/usr/share/xs-config/cfg/library/pgsql-xs/pg_hba.conf
