# plugins.d/startup/files/rc.local
#!/bin/bash
# this script is executed after all systemd startup has completed
#
xs-setup-network
