# plugins.d/startup/files/rc.local
#!/bin/bash
# this script is executed after all systemd startup has completed
#
# give NetworkManager time to actually configure the adapters
sleep 20
## JV this puts nameserver 127.0.0.1 in /etc/resolv.conf before 
## named is installed breaking dns
#xs-setup-network
# powerd generates its own rules wait a bit and load ours 
#systemctl condrestart iptables.service
#/bin/iiab-server &
