# roles/openvpn/files/rc.local
#!/bin/bash
#
# rc.local -- this file is executed at startup
