# plugins.d/idmgr/idmgr.manifest
/etc/systemd/system/idmgr.service.in
/etc/sysconfig/olpc-scripts/domain_config.d/idmgr
