# altfiles/etc/xinetd.d/xsactivation.in
# This is the configuration for the xs-activation server
# The service, as described by C. Scott Ananian
# * It listens on port 191 on those addresses on all interfaces.
# * When it receives a connection, it reads exactly 11 bytes.
# * It checks its "database" (text file, directory, in-memory hash) to
#   see if it has an activation lease for the serial number given by those
#   11 bytes.  If so, it sends it.
# * It closes the connection.


service xsactivation
{
	disable	        = no
	socket_type     = stream
	wait            = no
        port            = 191
	server          = /usr/bin/xs-activation-tcpserver.py

# Without type = UNLISTED xinetd wants this service to be called 'prospero' 
  	type		= UNLISTED	 

	#user		= root
	#group		= root
	user		= apache
	group		= apache

	# LAN IPv4
        flags           = IPv4
        bind            = 0.0.0.0

	# This service is unlikely to cause much resource
	# trouble, but rate-limit it anyway.
	instances  = 25
	max_load   = 5
	per_source = 5
	cps        = 150 10
	nice       = 10
	rlimit_as  = 75M
}
