# plugins.d/stats/stats.manifest
/etc/sugar-stats.conf.in
/etc/systemd/system/sugar-stats-server.service.in
